# files/etc/crontabs/root
# Use crontab.guru to set crontab schedule spesific time

0 * * * * /sbin/free.sh >/dev/null 2>&1
